# scripts/wakeup-mpd-display.service
[Unit]
Description=Wakeup MPD Display when playback starts

[Service]
ExecStart=/home/pi/wakeup-mpd-display.sh

[Install]
WantedBy=default.target
